# .docksal/services/cli/crontab
*/5 * * * * bash -l -c 'drush --root=/var/www cron'
